# 1m.service
[Unit]
Description=1m cron

[Service]
User=abraxas
ExecStart=/home/abraxas/myq/service/1m.sh
